# etc/supervisor/conf.d/ffmpeg.conf
[program:ffmpeg]
command=/usr/local/bin/ffmpeg.sh
user=selenium
priority=5
redirect_stderr=true
stdout_logfile_backups=2
stdout_logfile_maxbytes=10MB
autorestart=unexpected
